# deployment/systemd-reference/personal-website-db-backup-1.timer
[Unit]
Description=Timer for personal-website-db-backup-1

[Timer]
OnUnitActiveSec=4h
OnBootSec=1m
Persistent=true

[Install]
WantedBy=timers.target
